# com.upnp.mediaplayer/lib/application/scripts/systemd/mediaplayer.service
[Unit]
Description=MediaPlayer

[Service]
ExecStart=/home/pi/mediaplayer/run_systemd.sh

[Install]
WantedBy=multi-user.target